# systemctl_services/user/suspend.target
[Unit]
Description=User level suspend target
StopWhenUnneeded=yes
Wants=gnome-waydroid-rotator_user-restart.service
